# ai.carapace.gmail-proxy.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN"
  "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>ai.carapace.gmail-proxy</string>

    <key>ProgramArguments</key>
    <array>
        <string>/usr/local/bin/gmail-proxy</string>
        <string>serve</string>
        <string>--config</string>
        <string>/etc/carapace/gmail-proxy.toml</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>StandardOutPath</key>
    <string>/Users/carapace/.local/share/carapace/gmail-proxy.log</string>

    <key>StandardErrorPath</key>
    <string>/Users/carapace/.local/share/carapace/gmail-proxy.err</string>

    <key>EnvironmentVariables</key>
    <dict>
        <key>HOME</key>
        <string>/Users/carapace</string>
        <key>RUST_LOG</key>
        <string>info</string>
    </dict>
</dict>
</plist>
